# redis-statsd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="application/redis-statsd" type="service" version="0">
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>
    <exec_method name="start" type="method" exec="/opt/local/bin/python /opt/local/bin/redis-statsd.py --no-tags --period 10 --prefix {{ grains['host'] }}_ms_virtua_ch --redis-host 127.0.0.1:6379 --statsd-host stats.ms.virtua.ch:8125" timeout_seconds="60"/>
    <exec_method name="stop" type="method" exec=":kill" timeout_seconds="60"/>
    <instance name="default" enabled="true">
      <property_group name="startd" type="framework">
        <propval name="duration" type="astring" value="child"/>
        <propval name="ignore_error" type="astring" value="core,signal"/>
      </property_group>
      <property_group name="general" type="framework">
        <propval name="action_authorization" type="astring" value="solaris.smf.manage.redis-statsd"/>
        <propval name="value_authorization" type="astring" value="solaris.smf.value.redis-statsd"/>
      </property_group>
    </instance>
    <stability value="Evolving"/>
    <template>
      <common_name>
        <loctext xml:lang="Python">redis-statsd daemon</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
